# code/s6-overlay/s6-rc.d/celery/run
#!/bin/sh

celery -A app --workdir /app worker --loglevel=DEBUG  --concurrency=1